# server/lofi-stream-twitch.service
[Unit]
Description=Lofi Stream to Twitch
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/lofi-stream-twitch

# Set your Twitch stream key here
Environment=TWITCH_KEY=YOUR_TWITCH_STREAM_KEY_HERE

ExecStart=/opt/lofi-stream-twitch/stream.sh

# Auto-restart on crash
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
